# install/macos/com.snitch.menubar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.snitch.menubar</string>
    <key>ProgramArguments</key>
    <array>
        <string>__SNITCHBAR_APP__/Contents/MacOS/snitchbar</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/.snitch/snitchbar.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.snitch/snitchbar.log</string>
</dict>
</plist>
